# scripts/systemd/morning-news.timer
[Unit]
Description=早报推送定时器 - 早上9点

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
